# deploy/celery_worker.conf

# ==========================================
#  celery worker supervisor conf for Django
# ==========================================
# for django instance: spike.g-node.org
# generated by deploy.py - 2012-12-04 22:22:16.969788

[program:celery_worker]
command=/home/pmeier/Workspace/Gnode-Spike/gnode_spike/manage.py celery worker --loglevel=INFO
directory=/home/pmeier/Workspace/Gnode-Spike/gnode_spike
environment=PATH=""
user=www-data
numprocs=1
stdout_logfile=/home/pmeier/Workspace/Gnode-Spike/gnode_spike/log/celery_worker.log
stderr_logfile=/home/pmeier/Workspace/Gnode-Spike/gnode_spike/log/celery_worker.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; rabbitmq @ 1001, set this priority higher!
priority=1002

